# Packages/Controls/IQAudio/IQ_rot.service
[Unit]
Description=Starts the rotary encoder service
After=mpd.service

[Service]
PIDFile=/run/IQ_rot.pid
ExecStart=/usr/bin/iqaudio/IQ_rot
Restart=always

[Install]
WantedBy=multi-user.target
