# Packages/Controls/mute.service
[Unit]
Description=Starts the mute listener for MPD control
After=mpd.service

[Service]
PIDFile=/run/mute.pid
ExecStart=/usr/bin/controls/mute
Restart=always

[Install]
WantedBy=multi-user.target
