# Packages/Mongoose/mongoose.service
[Unit]
Description=Mongoose HTTP Daemon

[Service]
PIDFile=/run/mongoose.pid
ExecStart=/usr/bin/mongoose /etc/mongoose/mongoose.conf
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=always

[Install]
WantedBy=multi-user.target
